# deploy/laravel-worker.conf
; Example Supervisor config for C-by-AI Meal API workers
; Copy this file to /etc/supervisor/conf.d/laravel-worker.conf on the server,
; then run:
;   sudo supervisorctl reread
;   sudo supervisorctl update
;   sudo supervisorctl start laravel-worker
;
; Adjust paths, user, and numprocs as needed for your environment.

[program:laravel-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/24digi-AiMeal/artisan queue:work --sleep=3 --tries=3 --max-time=3600
autostart=true
autorestart=true
user=www-data
numprocs=4
redirect_stderr=true
stdout_logfile=/var/www/24digi-AiMeal/storage/logs/worker.log

